# distrib/build/sources/updates/update9/coretree.list
./etc/ppp/ip-up
./etc/logrotate.d/snort
./etc/rc.d/rc.firewall.up
./etc/snort.conf
