# distrib/build/sources/updates/update1/coretree.list
./etc/crontab
./etc/httpd/conf/httpd.conf
./etc/rc.d/rc.updatered
./httpd/cgi-bin/help.cgi
./httpd/cgi-bin/logs.cgi/log.dat
./usr/lib/smoothwall/services/apcupsd
