# distrib/build/sources/updates/update1/coretree.list
./etc/rc.d/restorescript
./etc/rc.d/rc.firewall.down
./etc/rc.d/rc.firewall.up
./etc/rc.d/rc.netaddress.up
./etc/rc.d/rc.sysinit
./httpd/cgi-bin/interfaces.cgi
./httpd/cgi-bin/updates.cgi
./usr/bin/smoothwall/updatelists.pl
./usr/bin/smoothwall/writeim.pl
./usr/lib/smoothwall/update.pm
