# distrib/build/sources/updates/update8/coretree.list
./etc/cdromextradevnodes
./etc/diskextradevnodes
./etc/fstab
./etc/nicdevices
./etc/rc.d/rc.sysinit
./etc/rc.d/rc.firewall.up
./sbin/hotplug
./usr/bin/installer/init
./usr/bin/installer/makeinitrd
./usr/bin/smoothwall/installgrub.pl
./usr/bin/smoothwall/updatestorageuuids.pl
./usr/bin/smoothwall/writefstab.pl
./usr/bin/smoothwall/writegrubconf.pl
./usr/etc/sshd_config
./usr/lib/smoothwall/fstab.boot
./usr/lib/smoothwall/fstab.in
./usr/lib/smoothwall/fstab.log
./usr/lib/smoothwall/fstab.root
./usr/lib/smoothwall/fstab.swap
./usr/lib/smoothwall/interface.pm
./usr/lib/smoothwall/menu.lst.in
./usr/lib/smoothwall/mergeparts
