# distrib/build/sources/updates/update2/coretree.list
./httpd/cgi-bin/ids.cgi
./httpd/cgi-bin/portfw.cgi
./etc/inittab
./etc/lilo.conf
./etc/cron.daily/runfreshclam
