# distrib/build/sources/updates/update6/coretree.list
./usr/lib/smoothwall/header.pm
./usr/lib/smoothwall/langs/alertboxes.base.pl
./usr/bin/installer/makeinitrd
./etc/nicdevices
./etc/modprobe.conf
./etc/lilo.conf
./etc/rc.d/rc.sysinit
./etc/rc.d/rc.firewall.up
./httpd/cgi-bin/ids.cgi
