# distrib/build/sources/updates/update7/coretree.list
./etc/rc.d/rc.adsl
./etc/rc.d/rc.netaddress.down
./etc/rc.d/rc.network
./httpd/cgi-bin/pppsetup.cgi
./usr/bin/smoothwall/writeproxy.pl
./usr/lib/smoothwall/langs/en.pl
./var/smoothwall/patches/installed
